# meal_display.service
[Unit]
Description=Meal Planner display daemon
After=network.target meal_planner.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/meal-planner-pi
ExecStart=/usr/bin/python3 /home/pi/meal-planner-pi/display.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
