# servelocal_planner.service
[Unit]
Description=ServeLocal web server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/var/lib/servelocal
ExecStart=/var/lib/servelocal/venv/bin/python3 /var/lib/servelocal/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
