# servelocal_display.service
[Unit]
Description=ServeLocal display daemon
# Deliberately does NOT wait on network.target or servelocal_planner.service:
# this should come up as early as possible so the screen shows a splash
# (proof of life) while the rest of boot -- network, the web server -- is
# still happening. display.py creates the DB schema itself if it doesn't
# exist yet, so it's safe regardless of which service wins the race.
After=local-fs.target
Wants=local-fs.target

[Service]
Type=simple
User=pi
WorkingDirectory=/var/lib/servelocal
ExecStart=/var/lib/servelocal/venv/bin/python3 /var/lib/servelocal/display.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
