# servelocal_planner.service
[Unit]
Description=ServeLocal web server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/servelocal
ExecStart=/usr/bin/python3 /home/pi/servelocal/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
